# admin/app.service
[Unit]
Description=Skymonitor Flask UI Service (reads from control API)
After=network.target control.service
Wants=control.service

[Service]
Type=simple
User=robert
Group=robert
WorkingDirectory=/home/robert/WeatherMonitor/safety-monitor
ExecStart=/home/robert/WeatherMonitor/venv/bin/python3 /home/robert/WeatherMonitor/safety-monitor/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# UI service can be stopped without affecting safety
# Control service remains independent

[Install]
WantedBy=multi-user.target
